# share/systemd/hermes-router-indicator.service.in
[Unit]
Description=Hermes Router Tray Indicator
After=network-online.target hermes-gateway.service
PartOf=hermes-gateway.service
StartLimitIntervalSec=0

[Service]
Type=simple
ExecStart=__EXEC_LINE__
Environment="HERMES_ROUTER_HOST=__HERMES_ROUTER_HOST__"
Restart=on-failure
RestartSec=5
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=210
# systemd --user services automatically inherit DISPLAY, DBUS_SESSION_BUS_ADDRESS,
# XDG_RUNTIME_DIR, etc. from the active graphical session — no manual Environment=
# injection is needed for GTK / AppIndicator to work.

[Install]
WantedBy=__WANTED_BY__
